# deployment/docker/supervisord_gunicorn.ini
[program:application]
command = gunicorn --config /tmp/gunicorn.conf.py siteapp.wsgi
directory = /usr/src/app
stderr_logfile = /var/log/application-stderr.log
stdout_logfile = /var/log/application-stdout.log

[program:notificationemails]
command = python3.6 manage.py send_notification_emails forever
directory = /usr/src/app
stderr_logfile = /var/log/notificationemails-stderr.log
stdout_logfile = /var/log/notificationemails-stdout.log
