# deployment/docker/supervisord.ini
[program:application]
command = gunicorn --config /tmp/gunicorn.conf.py siteapp.wsgi
directory = /usr/src/app
logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout

[program:notificationemails]
command = python3.6 manage.py send_notification_emails forever
directory = /usr/src/app
logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
