# deployment/docker/supervisord.ini
[program:application]
command = uwsgi --ini /tmp/uwsgi.ini
directory = /usr/src/app

[program:notificationemails]
command = python3.6 manage.py send_notification_emails forever
directory = /usr/src/app
